# videoTranscoding/root/etc/services.d/transcodeQueue/run
#!/bin/bash
exec s6-setuidgid abc /usr/local/bin/plexTranscode/queueman.sh
